# files/etc/uci-defaults/20_firewall_disable
#!/bin/sh
service firewall stop
service firewall disable
