# files/etc/rc.local
# This file is from the default image to store audio state

mkdir -p /var/lib/alsa
ln -s /etc/asound.state /var/lib/alsa/asound.state

exit 0
