# files/etc/rc.local
# This file is from the default image to store audio state

mkdir -p /etc/alsa
ln -s /etc/alsa/ /var/lib/

exit 0
